# packaging/debian/postinst
#!/bin/bash
set -e

# Create origindb user and group if they don't exist
if ! id "origindb" &>/dev/null; then
    useradd --system --shell /bin/false --home-dir /var/lib/origindb --create-home origindb
fi

# Create directories
mkdir -p /var/lib/origindb
mkdir -p /var/log/origindb
mkdir -p /etc/origindb

# Set ownership
chown -R origindb:origindb /var/lib/origindb
chown -R origindb:origindb /var/log/origindb
chown -R origindb:origindb /etc/origindb

# Set permissions
chmod 750 /var/lib/origindb
chmod 750 /var/log/origindb
chmod 755 /etc/origindb

# Create default configuration if it doesn't exist
if [ ! -f /etc/origindb/origindb.conf ]; then
    cat > /etc/origindb/origindb.conf << EOF
# OriginDB Configuration
# See https://origindb.com/docs/configuration for details

# Data directory
data_dir = /var/lib/origindb

# Network ports
websocket_port = 8080
grpc_port = 50051

# Logging
log_level = info
log_file = /var/log/origindb/origindb.log

# Security
enable_auth = false
EOF
    chown origindb:origindb /etc/origindb/origindb.conf
    chmod 640 /etc/origindb/origindb.conf
fi

# Enable and start systemd service
if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload
    systemctl enable origindb.service

    # Only start if this is a fresh install (not an upgrade)
    if [ "$1" = "configure" ] && [ -z "$2" ]; then
        echo "Starting OriginDB service..."
        systemctl start origindb.service
    fi
fi

echo "OriginDB has been installed successfully!"
echo "Configuration file: /etc/origindb/origindb.conf"
echo "Data directory: /var/lib/origindb"
echo "Log directory: /var/log/origindb"
echo ""
echo "To start OriginDB:"
echo "  sudo systemctl start origindb"
echo ""
echo "To enable auto-start on boot:"
echo "  sudo systemctl enable origindb"
echo ""
echo "For more information, visit: https://origindb.com/docs"