# packaging/debian/prerm
#!/bin/bash
set -e

# Stop and disable systemd service
if command -v systemctl >/dev/null 2>&1; then
    if systemctl is-active --quiet origindb.service; then
        echo "Stopping OriginDB service..."
        systemctl stop origindb.service
    fi

    if systemctl is-enabled --quiet origindb.service; then
        echo "Disabling OriginDB service..."
        systemctl disable origindb.service
    fi
fi

echo "OriginDB service stopped and disabled."
echo "Data in /var/lib/origindb has been preserved."
echo "To completely remove all data, run: sudo rm -rf /var/lib/origindb"