# packaging/linux/instantdb.service
[Unit]
Description=InstantDB Server
Documentation=https://instantdb.com/docs
After=network.target
Wants=network.target

[Service]
Type=simple
User=instantdb
Group=instantdb
ExecStart=/usr/bin/instantdb_server --data-dir /var/lib/instantdb --log-level info
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/instantdb /var/log/instantdb

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=INSTANTDB_DATA_DIR=/var/lib/instantdb
Environment=INSTANTDB_LOG_LEVEL=info
Environment=INSTANTDB_WS_PORT=8080
Environment=INSTANTDB_GRPC_PORT=50051

[Install]
WantedBy=multi-user.target